# dashboard.service
[Unit]
Description=Family Dashboard 9.7" ePaper Display
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/Work/FamilyDashboard97
Environment=PYTHONUNBUFFERED=1
Environment=BRING_EMAIL=your@email.com
Environment=BRING_PASSWORD=yourpassword
ExecStart=/usr/bin/python3 /home/pi/Work/FamilyDashboard97/main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
